# supervisor.conf
[program:space]
directory=/home/ubuntu/space/
command = /usr/bin/npm start
user = ubuntu
autostart=true
autorestart=true
stdout_logfile = /home/ubuntu/space/supervisor.log
redirect_stderr = true